# vncserverd/layout/DEBIAN/postinst
chown root:wheel /Library/LaunchDaemons/com.ethanarbuckle.vncserverd.plist
launchctl unload /Library/LaunchDaemons/com.ethanarbuckle.vncserverd.plist || true
launchctl load /Library/LaunchDaemons/com.ethanarbuckle.vncserverd.plist || true