# vncserverd/layout/Library/LaunchDaemons/com.ethanarbuckle.vncserverd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/vncserverd</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>UserName</key>
	<string>root</string>
    <key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.ethanarbuckle.vncserverd</string>
	<key>JetsamProperties</key>
	<dict>
		<key>JetsamPriority</key>
		<integer>18</integer>
	</dict>
</dict>
</plist>